# systemd/xorg@.socket
[Unit]
Description=Socket for xorg at display %i

[Socket]
ListenStrem=/tmp/.X11-unix/X%i
